# puppet/modules/topical/files/topical-worker.conf
start on starting topical
stop on stopping topical
